# comm_cpu-보드적용x_프로세스5개생성/image/참고/rc.local
#!/bin/bash
export PATH=$PATH:/root/comm

/root/comm/proc_main &
exit 0